# scripts/production/sanad-backend.service
[Unit]
Description=SANAD Platform Backend (Spring Boot)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=sanad
WorkingDirectory=/opt/sanad-platform
ExecStart=/usr/bin/java -XX:MaxRAMPercentage=75.0 -XX:+UseG1GC -Dfile.encoding=UTF-8 -jar /opt/sanad-platform/sanad-platform.jar --spring.profiles.active=prod
Restart=always
RestartSec=10
StandardOutput=append:/var/log/sanad-backend.log
StandardError=append:/var/log/sanad-backend.err
Environment=SPRING_PROFILES_ACTIVE=prod
Environment=DATABASE_URL=jdbc:postgresql://localhost:5432/sanad
Environment=DATABASE_USERNAME=sanad_backend
Environment=DATABASE_PASSWORD=Snad@2026!Backend
Environment=FLYWAY_ENABLED=true
Environment=JPA_DDL_AUTO=validate
Environment=SERVER_PORT=8080

[Install]
WantedBy=multi-user.target
